# fsw_allinone/support/fsw-node2.ini

[program:fsw-node2]
command=/bin/bash -c "/home/jboss/fsw/jboss-eap-6.1/bin/standalone.sh -c standalone-ha.xml -b 0.0.0.0 -bmanagement 0.0.0.0 -Djboss.server.base.dir=fsw2 -Djboss.node.name=fsw2 -Djboss.socket.binding.port-offset=1000"
user=jboss
