# fsw_allinone/support/httpd.ini

[program:httpd]
command=/bin/bash -c "exec /usr/sbin/apachectl -D FOREGROUND"
;command=/bin/bash -c "rm -rf /run/httpd/* && exec /usr/sbin/apachectl -D FOREGROUND"
;# Make sure we're not confused by old, incompletely-shutdown httpd
;# context after restarting the container.  httpd won't start correctly
;# if it thinks it is already running.